# bandwidth_occupier.service
[Unit]
Description=Oracle keep-alive bandwidth occupier
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
Environment=BANDWIDTH_SCRIPT=/usr/local/bin/bandwidth_occupier.sh
EnvironmentFile=-/etc/oalive/oalive.conf
ExecStart=/bin/sh ${BANDWIDTH_SCRIPT}
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
