# cpu-limit.service
[Unit]
Description=Oracle keep-alive CPU occupier
After=network-online.target

[Service]
Type=simple
User=root
Environment=CPU_SCRIPT=/usr/local/bin/cpu-limit.sh
EnvironmentFile=-/etc/oalive/oalive.conf
ExecStart=/bin/sh ${CPU_SCRIPT}
Restart=always
RestartSec=30s
KillMode=control-group

[Install]
WantedBy=multi-user.target
